# lib/ctfmod-load-secrets.service
[Unit]
Description=Load secrets into the ctfmod driver

[Service]
EnvironmentFile=/etc/default/ctfmod-flags
Type=oneshot
ExecStart=/usr/local/sbin/loader

[Install]
WantedBy=multi-user.target
